# packages/criscostack-brik/criscostack_brik-1.0.5-py3-none-any.whl/brik/config/templates/systemd/criscostack-brik-criscostack-schedule.service
[Unit]
Description="{{ brik_name }}-criscostack-schedule"
PartOf={{ brik_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ brik_cmd }} schedule
StandardOutput=file:{{ brik_dir }}/logs/schedule.log
StandardError=file:{{ brik_dir }}/logs/schedule.error.log
WorkingDirectory={{ brik_dir }}
