# packages/criscostack-brik/criscostack_brik-1.0.5-py3-none-any.whl/brik/config/templates/systemd/criscostack-brik-criscostack-long-worker.service
[Unit]
Description="{{ brik_name }}-criscostack-short-worker %I"
PartOf={{ brik_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ brik_cmd }} worker --queue long
StandardOutput=file:{{ brik_dir }}/logs/worker.log
StandardError=file:{{ brik_dir }}/logs/worker.error.log
WorkingDirectory={{ brik_dir }}
