# packages/criscostack-brik/criscostack_brik-1.0.5-py3-none-any.whl/brik/config/templates/systemd/criscostack-brik-web.target
[Unit]
After=network.target
Wants={{ brik_name }}-criscostack-web.service {{ brik_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
