# packages/criscostack-brik/criscostack_brik-1.0.5-py3-none-any.whl/brik/config/templates/systemd/criscostack-brik-redis.target
[Unit]
After=network.target
Wants={{ brik_name }}-redis-cache.service {{ brik_name }}-redis-queue.service

[Install]
WantedBy=multi-user.target
