# packages/criscostack-brik/criscostack_brik-1.0.5-py3-none-any.whl/brik/config/templates/systemd/criscostack-brik.target
[Unit]
After=network.target
Requires={{ brik_name }}-web.target {{ brik_name }}-workers.target {{ brik_name }}-redis.target

[Install]
WantedBy=multi-user.target
