# packages/criscostack-brik/criscostack_brik-1.0.5-py3-none-any.whl/brik/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:{{ brik_name }}-criscostack-web]
command={{ brik_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} --max-requests {{ gunicorn_max_requests }} --max-requests-jitter {{ gunicorn_max_requests_jitter }} -t {{ http_timeout }} criscostack.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/web.log
stderr_logfile={{ brik_dir }}/logs/web.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ brik_name }}-criscostack-schedule]
command={{ brik_cmd }} schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/schedule.log
stderr_logfile={{ brik_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ brik_dir }}

{% if not multi_queue_consumption %}
[program:{{ brik_name }}-criscostack-default-worker]
command={{ brik_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/worker.log
stderr_logfile={{ brik_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ brik_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
{% endif %}

[program:{{ brik_name }}-criscostack-short-worker]
command={{ brik_cmd }} worker --queue short{{',default' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/worker.log
stderr_logfile={{ brik_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ brik_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ brik_name }}-criscostack-long-worker]
command={{ brik_cmd }} worker --queue long{{',default,short' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/worker.log
stderr_logfile={{ brik_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ brik_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

{% for worker_name, worker_details in workers.items() %}
[program:{{ brik_name }}-criscostack-{{ worker_name }}-worker]
command={{ brik_cmd }} worker --queue {{ worker_name }}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/worker.log
stderr_logfile={{ brik_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs={{ worker_details["timeout"] }}
directory={{ brik_dir }}
killasgroup=true
numprocs={{ worker_details["background_workers"] or background_workers }}
process_name=%(program_name)s-%(process_num)d
{% endfor %}


{% if not skip_redis %}
[program:{{ brik_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/redis-cache.log
stderr_logfile={{ brik_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ brik_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/redis-queue.log
stderr_logfile={{ brik_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if node %}
[program:{{ brik_name }}-node-socketio]
command={{ node }} {{ brik_dir }}/apps/criscostack/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ brik_dir }}/logs/node-socketio.log
stderr_logfile={{ brik_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ brik_dir }}
{% endif %}

[group:{{ brik_name }}-web]
programs={{ brik_name }}-criscostack-web {%- if node -%} ,{{ brik_name }}-node-socketio {%- endif%}


{% if multi_queue_consumption %}

[group:{{ brik_name }}-workers]
programs={{ brik_name }}-criscostack-schedule,{{ brik_name }}-criscostack-short-worker,{{ brik_name }}-criscostack-long-worker{%- for worker_name in workers -%},{{ brik_name }}-criscostack-{{ worker_name }}-worker{%- endfor %}

{% else %}

[group:{{ brik_name }}-workers]
programs={{ brik_name }}-criscostack-schedule,{{ brik_name }}-criscostack-default-worker,{{ brik_name }}-criscostack-short-worker,{{ brik_name }}-criscostack-long-worker{%- for worker_name in workers -%},{{ brik_name }}-criscostack-{{ worker_name }}-worker{%- endfor %}

{% endif %}

{% if not skip_redis %}
[group:{{ brik_name }}-redis]
programs={{ brik_name }}-redis-cache,{{ brik_name }}-redis-queue
{% endif %}
